# nullmailer/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/nullmailer
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwsr-xr-x	nullmail/root	usr/bin/mailq
-rwxr-xr-x	root/root	usr/bin/nullmailer-dsn
-rwxr-xr-x	root/root	usr/bin/nullmailer-inject
-rwxr-xr-x	root/root	usr/bin/nullmailer-smtpd
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/nullmailer/
-rw-r--r--	root/root	usr/etc/nullmailer/me
-rw-r--r--	root/root	usr/etc/nullmailer/remotes (EMPTY)
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/nullmailer/
drwxr-xr-x	root/root	usr/lib/nullmailer/nullmailer/
-rwxr-xr-x	root/root	usr/lib/nullmailer/nullmailer/qmqp
-rwxr-xr-x	root/root	usr/lib/nullmailer/nullmailer/smtp
drwxr-xr-x	root/root	usr/sbin/
-rwsr-xr-x	nullmail/root	usr/sbin/nullmailer-queue
-rwxr-xr-x	root/root	usr/sbin/nullmailer-send
-rwxr-xr-x	root/root	usr/sbin/sendmail
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/nullmailer-dsn.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nullmailer-inject.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sendmail.1.gz
drwxr-xr-x	root/root	usr/share/man/man7/
-rw-r--r--	root/root	usr/share/man/man7/nullmailer.7.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/nullmailer-queue.8.gz
-rw-r--r--	root/root	usr/share/man/man8/nullmailer-send.8.gz
drwxr-xr-x	root/root	usr/var/
drwxr-xr-x	root/root	usr/var/spool/
drwxr-xr-x	root/root	usr/var/spool/nullmailer/
drwx------	nullmail/root	usr/var/spool/nullmailer/failed/
drwx------	nullmail/root	usr/var/spool/nullmailer/queue/
drwx------	nullmail/root	usr/var/spool/nullmailer/tmp/
prw-------	nullmail/root	usr/var/spool/nullmailer/trigger
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
-rw-r--r--	nullmail/daemon	var/log/nullmail (EMPTY)
